# systemd/fcgiwrap.socket.d/mode.conf
[Socket]
SocketMode=0660
SocketUser=the-tangle
SocketGroup=www-data